# mkbackup-btrfs_0.2_all/etc/systemd/system/timer-minute.timer
[Unit]
Description=1-Minute Timer

[Timer]
OnUnitActiveSec=1min
Unit=timer-minute.target

[Install]
WantedBy=timers.target
