# mkbackup-btrfs/etc/systemd/system/backup.path
[Unit]
Description=Starts mounting backups to /backup, if this directory exists
BindsTo=mkbackup.target

[Path]
PathExists=/backup
Unit=backup.automount

[Install]
WantedBy=mkbackup.target
