# mkbackup-btrfs/etc/systemd/system/backup.mount
[Unit]
Before=local-fs.target

[Mount]
What=/etc/systemd/system/scripts/MksnapshotFS.py
Where=/backup
Type=fuse
TimeoutSec=10s
Options=ro,nofail
