# mkbackup-btrfs/etc/systemd/system/mkbackup-userdir@.path
[Unit]
Description=Starts mounting backups to /home/%i/backup, if this directory exists
BindsTo=mkbackup-userdir.target

[Path]
PathExists=/home/%i/backup
Unit=mkbackup-userdir@%i.service

[Install]
WantedBy=mkbackup-userdir.target
