# mkbackup-btrfs/etc/systemd/system/timer-aptupgrade.timer
[Unit]
Description=Runs backup %I before system-update/upgrade with apt
#BindsTo=mkbackup@BKP.target
PartOf=mkbackup@SNP.target

[Timer]
OnActiveSec=1s
AccuracySec=2s
Unit=timer-aptupgrade.target
RemainAfterElapse=false

[Install]
WantedBy=mkbackup@SNP.target
