# mkbackup-btrfs/usr/lib/systemd/user/mkbackup-userdir.path
[Unit]
Description=Starts mounting backups to %h/backup, if this directory exists
#BindsTo=mkbackup-userdir.path
Conflicts=shutdown.target

[Path]
PathExists=%h/backup
Unit=mkbackup-userdir.service

[Install]
WantedBy=paths.target
