# mkbackup-btrfs/usr/lib/systemd/user/mkbackup-userdir.service
[Unit]
Description=Show all userspezific backups in %h/backup
BindsTo=mkbackup-userdir.path
Conflicts=shutdown.target

[Service]
ExecStart=/usr/bin/MksnapshotFS.py -f %h/backup -o ro,allow_root
ExecStop=-/bin/fusermount -u %h/backup
Restart=on-success
