# etc/systemd/user/cups-fuse.path
[Unit]
Description=Starts mounting printers to %h/Desktop/printers, if this directory exists

[Path]
PathExists=%h/Desktop/Drucker
Unit=cups-fuse.service

[Install]
WantedBy=paths.target
