# etc/systemd/user/mkbackup-userdir.path
[Unit]
Description=Starts mounting backups to %h/backup, if this directory exists

[Path]
PathExists=%h/backup
Unit=mkbackup-userdir.service

[Install]
WantedBy=paths.target
