# systemd-units/timer-afterboot.target
[Unit]
Description=Target after successfull bootup
StopWhenUnneeded=yes
