# mkbackup-btrfs/lib/systemd/system/backup.path
[Unit]
Description=Start automounter for /backup, if this directory exists
PartOf=mkbackup.target
#Before=local-fs.target
#DefaultDependencies = no

[Path]
PathExists=/backup
Unit=backup.automount

[Install]
WantedBy=mkbackup.target
