# mkbackup-btrfs/lib/systemd/system/backup.automount
[Unit]
Description=Automounting of /backup - activated by backup.path
Before=local-fs.target
PartOf=backup.path

[Automount]
TimeoutIdleSec=15s
Where=/backup
