# mkbackup-btrfs/lib/systemd/system/mkbackup-conf@.path
[Unit]
Description=Upate /tmp/%I.conf.tmp if /etc/%I* changes

[Path]
PathModified=/etc/%I.conf
PathModified=/etc/%I.conf.d/
Unit=mkbackup-conf@%i.service

[Install]
WantedBy=paths.target
