# mkbackup-btrfs/lib/systemd/system/smartctl-fast@.service
[Unit]
Description=Smartctl Health-Test for %f

[Service]
ExecStart=/usr/sbin/smartctl -H %f
