# mkbackup-btrfs/lib/systemd/system/timer-plugin.target
[Unit]
Description=Target after successfully plugging in external backup-drive
StopWhenUnneeded=yes
