# mkbackup-btrfs/lib/systemd/system/timer-aptupgrade.target
[Unit]
Description=Triggered by apt
StopWhenUnneeded=yes
