# etc/systemd/system/backup.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
PartOf=var-cache-backup.mount

[Mount]
What=/var/cache/backup/aldebaran
Where=/backup
Type=none
Options=bind,ro,nofail,noauto,comment=systemd.automount
