# etc/systemd/system/backup.automount
[Unit]
Before=local-fs.target
BindsTo=backup-automount.target

[Automount]
Where=/backup

[Install]
WantedBy=backup-automount.target
