# etc/systemd/system/mkbackup-userdir.target
[Unit]
Description=Enables mounting of all ~/backup directories


[Install]
WantedBy=mkbackup.target
