# etc/systemd/system/mkbackup@BKP.target
[Unit]
Description=Target for services for %i in mkbackup

Documentation=man:systemd.special(7)
StopWhenUnneeded=no

