# etc/systemd/system/mkbackup.target
[Unit]
Description=Enables all backup.services


[Install]
WantedBy=basic.target
