# etc/systemd/system/timer-manually.target
[Unit]
Description=Target for manually backups
StopWhenUnneeded=yes
