# etc/systemd/system/timer-minute.timer
[Unit]
Description=1-Minute Timer

[Timer]
OnUnitActiveSec=1min
Unit=timer-dmin.target

[Install]
WantedBy=timers.target
