# etc/systemd/system/timer-minute.target
[Unit]
Description=Every 1 Minutes Timer Target
StopWhenUnneeded=yes
