# etc/systemd/system/timer-dmin.target
[Unit]
Description=Every 10 Minutes Timer Target
StopWhenUnneeded=yes
